# systemd/poll-camera.service
[Unit]
After=graphical.service

[Service]
ExecStart=/usr/local/sbin/poll_camera.sh

[Install]
WantedBy=default.target
